# deploy/relay.service
[Unit]
Description=Relay WebSocket signaling service
After=network-online.target coturn.service
Wants=network-online.target coturn.service

[Service]
Type=simple
User=relay
Group=relay
WorkingDirectory=/opt/relay
EnvironmentFile=/etc/relay/relay.env
ExecStart=/usr/bin/node /opt/relay/index.js
Restart=always
RestartSec=3
TimeoutStopSec=15
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
LockPersonality=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
